# .config/systemd/user/grub-notify.service
[Unit]
Description=Notify when GRUB config is modified

[Service]
Type=oneshot
ExecStart=%h/.local/bin/grub-notify.sh
